# src/argon-dashboard.service
[Unit]
Description=Argon ONE UP CM5 Dashboard Daemon
After=multi-user.target
Requires=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /usr/local/bin/argon_daemon.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argon-dashboard

# I2C Zugriff
SupplementaryGroups=i2c

# Sicherheit - gelockert fuer /sys/class/hwmon und /sys/class/leds Zugriff
ProtectSystem=full
ReadWritePaths=/tmp /sys/class/hwmon /sys/class/leds
DeviceAllow=/dev/i2c-1 rw
ProtectHome=true
ProtectKernelModules=true

[Install]
WantedBy=multi-user.target
